# .hdd/bitseed_ui.conf
#start bitseed GUI settings daemon script
start on startup
task
# exec su - linaro -c "/home/linaro/ui_init.sh"
script
   su - linaro -c "/home/linaro/ui_init.sh"
end script
